# honeypot/config/etc/systemd/system/dnsmasq.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target systemd-networkd.service

[Service]
Restart=on-failure
RestartSec=5s
